# scripts/pc-stats-agent.service
[Unit]
Description=Worker Monitoring Agent
After=network.target

[Service]
Type=simple
User=nobody
WorkingDirectory=/opt/worker-agent
ExecStart=/opt/worker-agent/agent
Restart=on-failure
RestartSec=10

# Environment variables
Environment="PORT=8080"
Environment="INTERVAL_SEC=30"
Environment="HISTORY_SIZE=720"

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=worker-agent

[Install]
WantedBy=multi-user.target
